# docker/conf/api.conf
[program:royal]
directory=/var/www/html/

command=node dist/src/main.js
user=root
autostart=true
autorestart=true
redirect_stderr=true
stopasgroup=true
killasgroup=true
